# src/data/services.yml
name: "Windows Services"
description: "Registry keys defining Windows services configuration including binary paths, start types, and dependencies. Commonly abused for persistence by creating malicious services or modifying existing ones."
tags:
  - persistence
  - privilege_escalation
  - execution
mitre_ids:
  - T1543.003
keys:
  - 'HKLM\SYSTEM\CurrentControlSet\Services'
  - 'HKLM\SYSTEM\CurrentControlSet\Enum\Root\LEGACY_*'
  - 'HKLM\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Svchost'
resources:
  - link: https://learn.microsoft.com/en-us/windows/win32/services/services
  - link: https://attack.mitre.org/techniques/T1543/003/
use_cases:
  - title: "Create a malicious service"
    description: "Creates a new Windows service that runs a payload with SYSTEM privileges."
    code: |
      sc create MaliciousSvc binPath= "C:\Windows\Temp\payload.exe" start= auto
  - title: "Modify existing service binary path"
    description: "Hijacks an existing service by changing its binary path."
    code: |
      reg add "HKLM\SYSTEM\CurrentControlSet\Services\TargetSvc" /v "ImagePath" /t REG_EXPAND_SZ /d "C:\Windows\Temp\payload.exe" /f
  - title: "Query all services"
    description: "Lists all registered services with their configurations."
    code: |
      reg query "HKLM\SYSTEM\CurrentControlSet\Services" /s | findstr /i "ImagePath"
  - title: "Detect svchost-hosted services"
    description: "Lists services that run inside svchost.exe processes."
    code: |
      reg query "HKLM\SOFTWARE\Microsoft\Windows NT\CurrentVersion\Svchost"
